# templates/mailcatcher/mailcatcher.service.j2
[Unit]
Description=MailCatcher Service
After=network.service vagrant.mount
[Service]
Type=simple
ExecStart=/usr/local/bin/mailcatcher --foreground
User=root
[Install]
WantedBy=multi-user.target
